# .docker/src-cms/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P

/usr/sbin/nginx -g "error_log /proc/self/fd/2 error;" -c /etc/nginx/nginx.conf

# s6-setuidgid nginx
