# .docker/src-waf/etc/periodic/hourly/update-nginx-deny-by-domain.sh
#!/bin/sh
/etc/nginx/dynamic/deny_domain_by_ip.sh

