# .docker/src-cms/etc/services.d/cron/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/crond -f -c /etc/crontabs
