# .docker/src-cms/etc/services.d/cron/run
#!/usr/bin/with-contenv sh

exec cron -f
