# .docker/src-cron/etc/services.d/cert-watcher/run
#!/usr/bin/with-contenv sh

exec /opt/cron/cert-watcher.sh
