# .docker/src-cms/etc/s6-overlay/s6-rc.d/appserver-prepare/up
/command/with-contenv /etc/s6-overlay/scripts/appserver-prepare
