# .docker/src-cms/etc/s6-overlay/s6-rc.d/cron/run
#!/command/execlineb -P

/usr/sbin/crond -f -c /etc/crontabs
