# .docker/src-cms/etc/s6-overlay/s6-rc.d/cron/run
#!/command/execlineb -P
with-contenv
/usr/sbin/crond -f -c /etc/crontabs
