# .docker/src-cms/etc/s6-overlay/s6-rc.d/web/run
#!/command/execlineb -P

/usr/sbin/nginx -g "daemon off;"
