# .docker/src-cms/etc/s6-overlay/s6-rc.d/webserver-prepare/up
/command/with-contenv /etc/s6-overlay/scripts/webserver-prepare
