# .docker/src-www/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

exec nginx -g "error_log /dev/stderr error;"
