# .docker/src-cms/etc/services.d/beta/run
#!/usr/bin/with-contenv sh
www=/var/www
sleep 300
. ${www}/scripts/beta/betapush.sh