# .docker/src-cms/etc/s6-overlay/s6-rc.d/bootstrap/up
#!/usr/bin/execlineb -P

/var/www/scripts/bootstrap.sh
