# .docker/src-cms/etc/services.d/cert-watcher/run
#!/usr/bin/with-contenv sh

exec /var/www/scripts/cert-watcher.sh
