# .buildpack/src-cronrunner/crontab
*    *       *       *       *       run-parts ~/app/etc/periodic/1min
