# .docker/src-cms/etc/periodic/15min/generate-static-site

# need logic here to grab name of env and pass in to tome-run the correct uri for this environment

# only the 1st instance within cloud formation should actually do anything on cron
if [ "${CF_INSTANCE_INDEX:-''}" == "0" ]; then
  /var/www/scripts/tome-run.sh
fi;
