# .docker/src-cms/etc/periodic/15min/generate-static-site

if [ "${CF_INSTANCE_INDEX:-''}" == "0" ]; then
  /var/www/scripts/tome-run.sh
fi;
