# .docker/src-cms/etc/s6-overlay/s6-rc.d/appserver/run
#!/command/execlineb -P
with-contenv
fdmove -c 2 1
php-fpm -F
