# .docker/src-cms/etc/s6-overlay/s6-rc.d/php/run
#!/command/execlineb -P

/usr/sbin/php-fpm -F

# s6-setuidgid nginx
